# scripts/systemd/a-term-backend.service
[Unit]
Description=A-Term Backend (FastAPI)
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_ROOT__
Environment="PATH=%h/bin:%h/.local/bin:__PROJECT_ROOT__/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=%h"
EnvironmentFile=-__PROJECT_ROOT__/.env.local
EnvironmentFile=-__PROJECT_ROOT__/.env
ExecStart=/bin/bash -lc 'export A_TERM_BIND_HOST="${A_TERM_BIND_HOST:-127.0.0.1}"; export A_TERM_PORT="${A_TERM_PORT:-8002}"; exec __PROJECT_ROOT__/.venv/bin/uvicorn a_term.main:app --host "$A_TERM_BIND_HOST" --port "$A_TERM_PORT"'
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=a-term-backend

[Install]
WantedBy=default.target
