# scripts/systemd/a-term-frontend.service
[Unit]
Description=A-Term Frontend (Next.js)
After=network.target a-term-backend.service
Wants=a-term-backend.service

[Service]
Type=simple
WorkingDirectory=__PROJECT_ROOT__/frontend
Environment="PATH=%h/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=%h"
EnvironmentFile=-__PROJECT_ROOT__/.env.local
EnvironmentFile=-__PROJECT_ROOT__/.env
Environment="NODE_ENV=production"
ExecStart=/bin/bash -lc 'export API_URL="${API_URL:-http://127.0.0.1:${A_TERM_PORT:-8002}}"; export HOSTNAME="${A_TERM_FRONTEND_HOST:-127.0.0.1}"; export PORT="${A_TERM_FRONTEND_PORT:-3002}"; exec corepack pnpm start --hostname "$HOSTNAME" --port "$PORT"'
TimeoutStopSec=10
KillMode=mixed
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=a-term-frontend

[Install]
WantedBy=default.target
